# net-vpn/pritunl-client-electron-service/files/pritunl-client.service
[Unit]
Description=Pritunl Client Daemon

[Service]
ExecStart=/usr/bin/pritunld

[Install]
WantedBy=multi-user.target
